# init/runit/log/run
#!/bin/sh
exec svlogd -tt @@INSTALL_DIR@@/log
