# init/systemd/gzw-eac-fix.path
[Unit]
Description=Watch for Gray Zone Warfare updates

[Path]
PathChanged=@@MANIFEST_PATH@@

[Install]
WantedBy=default.target
